# files/system/etc/profile.d/just.sh
alias just='just --justfile /usr/share/ublue-os/just/60-custom.just'
